# deploy/clubfridge-kasse.service
[Unit]
Description=Clubfridge Kasse
Documentation=https://github.com/clubfridge-dev/clubfridgeV2
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
User=%i
# Service-Mode rendert direkt über KMSDRM (Drop-in kmsdrm.conf setzt
# SDL_VIDEODRIVER=kmsdrm). Wir setzen daher KEIN DISPLAY/XAUTHORITY —
# sonst fällt SDL2 bei blockiertem KMSDRM stillschweigend auf Xwayland
# zurück und die Rotation wird gegenüber dem Compositor verdoppelt.
Environment=KIVY_NO_ENV_CONFIG=1

WorkingDirectory=/opt/clubfridge/kasse
ExecStart=/opt/clubfridge/kasse/.venv/bin/python main.py

Restart=always
RestartSec=5
# Prozess darf max. 3x innerhalb von 60s neu starten
# danach wartet systemd länger (StartLimitBurst=3)
StartLimitBurst=3
StartLimitIntervalSec=60

StandardOutput=journal
StandardError=journal
SyslogIdentifier=clubfridge-kasse

# Sandboxing deaktiviert: WLAN-Setup braucht sudo/nmcli,
# und mehrere Protect*-Optionen setzen implizit NoNewPrivileges=true

[Install]
WantedBy=graphical.target multi-user.target
