# deploy/clubfridge-update@.timer
[Unit]
Description=Clubfridge Kasse – Täglicher Update-Check (%i)
Documentation=https://github.com/tb59427/clubfridge-kasse

[Timer]
# Täglich um 03:00 Uhr
OnCalendar=*-*-* 03:00:00
# Zufällige Verzögerung bis zu 10 Minuten (verhindert gleichzeitige Updates bei mehreren Kassen)
RandomizedDelaySec=600
# Verpasster Lauf (z. B. Pi war aus) wird beim nächsten Start nachgeholt
Persistent=true

[Install]
WantedBy=timers.target
